# src/ansible/templates/etc/systemd/system/archive.timer
[Unit]
Description=Run archive.service weekly

[Timer]
OnCalendar=Sat *-*-* 00:00:00 UTC

[Install]
WantedBy=timers.target