# src/ansible/templates/etc/systemd/system/archive.service
[Unit]
Description=Rotate log files to /var/archive

[Service]
ExecStart=/etc/systemd/system/mail-run "%H %n" /etc/systemd/system/archive