# src/ansible/templates/etc/systemd/system/cvsd.service
[Unit]
Description=cvsd

[Service]
Type=forking
ExecStart=/usr/sbin/cvsd
PIDFile=/var/run/cvsd.pid

[Install]
WantedBy=multi-user.target
